# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:mariadb]
command=/usr/bin/mysqld_safe

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
